# docs/systemd/open-tag-server.service
# open-tag control-plane server — systemd unit
#
# SETUP (run as root or with sudo):
#   1. Edit the [Service] section below:
#      - Replace /home/ubuntu/open-tag with the actual checkout path
#      - Replace ubuntu:ubuntu with the actual user:group running the server
#   2. Copy to systemd:
#        sudo cp docs/systemd/open-tag-server.service /etc/systemd/system/open-tag-server.service
#   3. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now open-tag-server
#   4. Check status:
#        sudo systemctl status open-tag-server
#        journalctl -u open-tag-server -f
#
# Prerequisites: Node.js 20+ on PATH, .env.prod configured, web/dist built.
# See docs/self-host.md "Bare Node.js (advanced)" for the full procedure.

[Unit]
Description=open-tag control-plane server
Documentation=https://github.com/fancyboi999/open-tag
After=network.target postgresql.service redis.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/open-tag

# Load .env.prod before starting (Node's native loadEnvFile reads the file
# set in ENV_FILE; the ExecStartPre step migrates the DB schema first).
Environment=ENV_FILE=.env.prod
Environment=NODE_ENV=production

# Migrate schema before starting the server (additive-safe; fails loud on
# destructive changes, which will prevent the server from starting — intentional).
ExecStartPre=/usr/bin/env bash -c \
  "set -a; . /home/ubuntu/open-tag/.env.prod; set +a; npx drizzle-kit push"

ExecStart=/usr/bin/env npx tsx src/server/index.ts

Restart=on-failure
RestartSec=5s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=open-tag-server

# Graceful shutdown: allow up to 30s for in-flight requests.
TimeoutStopSec=30

[Install]
WantedBy=multi-user.target
