# docs/systemd/open-tag-daemon.service
# open-tag compute-plane daemon — systemd unit
#
# This unit is for the BARE NODE.JS path (see docs/self-host.md "Bare Node.js").
# It requires: the repo cloned at WorkingDirectory, and tsx available on PATH.
#
# DOCKER COMPOSE users: the daemon still runs on the HOST (not in a container),
# but you should use the published npm package form (no repo needed):
#   ExecStart=/usr/bin/env npx @fancyboi999/open-tag-daemon@latest \
#     --server-url http://localhost:7788 \
#     --api-key change-me
# (7788 is the Docker Compose default port; adjust OPEN_TAG_SERVER_URL below accordingly.)
#
# SETUP (run as root or with sudo):
#   1. Edit the [Service] section below:
#      - Set WorkingDirectory to the actual checkout path
#      - Set User/Group to the user running the daemon
#      - Set OPEN_TAG_SERVER_URL to where the server is reachable
#        (bare Node.js default: http://localhost:7777;
#         Docker Compose default: http://localhost:7788)
#      - Set OPEN_TAG_API_KEY to match DAEMON_BOOTSTRAP_KEY in your .env.prod / .env.docker
#   2. Copy to systemd:
#        sudo cp docs/systemd/open-tag-daemon.service /etc/systemd/system/open-tag-daemon.service
#   3. Enable and start:
#        sudo systemctl daemon-reload
#        sudo systemctl enable --now open-tag-daemon
#   4. Check status:
#        sudo systemctl status open-tag-daemon
#        journalctl -u open-tag-daemon -f
#
# NOTE on nvm / PATH: systemd runs with a minimal PATH that may not include nvm-managed
# Node.js. If `npx` is not found, use the full absolute path (e.g.
# /home/ubuntu/.nvm/versions/node/v22.x.x/bin/npx) or add an Environment=PATH= line.
# See docs/self-host.md "Start with systemd" for the full explanation.
#
# See docs/self-host.md for the complete procedure.

[Unit]
Description=open-tag compute-plane daemon
Documentation=https://github.com/fancyboi999/open-tag
After=network.target
# If the server runs on the same host (bare Node.js path), wait for it:
After=open-tag-server.service
Wants=network.target

[Service]
Type=simple
User=ubuntu
Group=ubuntu
WorkingDirectory=/home/ubuntu/open-tag

# Server URL — set to where the control-plane server is reachable.
# Bare Node.js default: http://localhost:7777
# Docker Compose default: http://localhost:7788
Environment=OPEN_TAG_SERVER_URL=http://localhost:7777

# Must match DAEMON_BOOTSTRAP_KEY in the server's .env.prod / .env.docker.
# Replace with the real value, or load from a secrets file:
#   EnvironmentFile=/etc/open-tag/daemon.env
Environment=OPEN_TAG_API_KEY=change-me

# OPEN_TAG_HOME: where agent workspaces, MEMORY.md files, logs, and uploads are stored.
# Defaults to ~/.open-tag (i.e. /home/ubuntu/.open-tag for User=ubuntu).
# Override here if you need a different path:
# Environment=OPEN_TAG_HOME=/data/open-tag

# Bare Node.js form (requires repo checkout + tsx in devDependencies):
ExecStart=/usr/bin/env npx tsx src/daemon/index.ts \
  --server-url ${OPEN_TAG_SERVER_URL} \
  --api-key ${OPEN_TAG_API_KEY}

# Docker Compose / no-repo form — uncomment and comment out the bare form above:
# ExecStart=/usr/bin/env npx @fancyboi999/open-tag-daemon@latest \
#   --server-url ${OPEN_TAG_SERVER_URL} \
#   --api-key ${OPEN_TAG_API_KEY}

Restart=on-failure
RestartSec=10s
StandardOutput=journal
StandardError=journal
SyslogIdentifier=open-tag-daemon

[Install]
WantedBy=multi-user.target
